# supervisor_automysqlbackup.conf
[program:automysqlbackup]
command=/usr/sbin/automysqlbackup -c /etc/default/automysqlbackup

