# supervisor_automysqlbackup.conf
[program:automysqlbackup]
/usr/sbin/automysqlbackup -c /etc/default/automysqlbackup

